# deploy/examples/supervisor-teamhub.conf
; TeamHUB — Supervisor configuration (example)
;
; Install: sudo cp deploy/examples/supervisor-teamhub.conf /etc/supervisor/conf.d/teamhub.conf
;          sudo supervisorctl reread && sudo supervisorctl update
;
; Replace /var/www/teamhub and www-data with your deploy path and user.

[program:teamhub-octane]
process_name=%(program_name)s
command=php /var/www/teamhub/artisan octane:start --server=roadrunner --host=127.0.0.1 --port=8000
directory=/var/www/teamhub
user=www-data
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
redirect_stderr=true
stdout_logfile=/var/log/supervisor/teamhub-octane.log
stopwaitsecs=30

[program:teamhub-queue]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/teamhub/artisan queue:work redis --sleep=3 --tries=3 --max-time=3600
directory=/var/www/teamhub
user=www-data
numprocs=1
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
redirect_stderr=true
stdout_logfile=/var/log/supervisor/teamhub-queue.log
stopwaitsecs=60
